# djangoapp/templates/supervisor-django.conf.in
[program:django-uwsgi]
command = ${uwsgi_binary} --ini ${uwsgi_config}

[program:nginx-django]
command = /usr/sbin/nginx
